# testdata/generated/k-http.target
# generated by k
[Unit]
Requires=k-http.socket k-https.socket k-http.service

