# testdata/generated/k-http.service
# generated by k
[Service]
AmbientCapabilities=CAP_NET_BIND_SERVICE
CacheDirectory=k
CapabilityBoundingSet=CAP_NET_BIND_SERVICE
DynamicUser=true
Environment=K_RUN_DIR=testdata/generated/k
ExecReload=kill -USR1 $MAINPID
ExecStart=k serve ${K_RUN_DIR}/server.json
LogExtraFields=K=k-http
StateDirectory=k-http
SyslogIdentifier=k-http

[Unit]
Description=k web server
PartOf=k-http.target

