# testdata/generated/k.target
# generated by k
[Unit]
After=network-online.target
Requires=example_app.target k-http.target

