# config/testdata/generated/app.target
# generated by k
[Unit]
OnFailure=k-notify@%N.service
Requires=app.service

