# config/testdata/generated/app.service
# generated by k
[Service]
CacheDirectory=app
DynamicUser=true
Environment=K_RUN_DIR=testdata/tmp/k
Environment=foo=bar
EnvironmentFile=testdata/tmp/k/app.env
ExecStart=/opt/k/app/main
LogExtraFields=K=app
Restart=always
StateDirectory=app
SyslogIdentifier=app
WorkingDirectory=/var/lib/app/

[Unit]
PartOf=app.target k.target

