# config/testdata/generated/k-http.target
# generated by k
[Unit]
OnFailure=k-notify@%N.service
Requires=k-http.service k-http.socket k-https.socket

