# config/testdata/generated/k-http.service
# generated by k
[Service]
CacheDirectory=k-http
DynamicUser=true
Environment=K_RUN_DIR=testdata/tmp/k
EnvironmentFile=testdata/tmp/k/k-http.env
ExecStart=/usr/bin/echo serve ${K_RUN_DIR}/k-http.json
LogExtraFields=K=k-http
Restart=always
StateDirectory=k-http
SyslogIdentifier=k-http

[Unit]
PartOf=k-http.target k.target

