# config/testdata/generated/k-https.socket
# generated by k
[Socket]
FileDescriptorName=https
ListenStream=443
Service=k-http.service

[Unit]
PartOf=k-http.target k.target

