# config/testdata/generated/k-http.socket
# generated by k
[Socket]
FileDescriptorName=http
ListenStream=80
Service=k-http.service

[Unit]
PartOf=k-http.target k.target

