# config/testdata/generated/k-notify@.service
# generated by k
[Service]
ExecStart=-/usr/bin/echo notify --app %i

