# config/testdata/generated/k.target
# generated by k
[Unit]
After=network-online.target
Requires=app.target k-http.target

